# scripts/systemd/zai-quota.service
[Unit]
Description=ZAI Usage Quota Monitor Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
Environment="PATH={{PROJECT_ROOT}}/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile={{PROJECT_ROOT}}/.env

ExecStart={{PROJECT_ROOT}}/.venv/bin/python {{PROJECT_ROOT}}/src/main.py --server --port ${PORT:-9999}

StandardOutput=journal
StandardError=journal
SyslogIdentifier=zai-quota

Restart=always
RestartSec=10

[Install]
WantedBy=default.target
